# data/init/openrc/glogikd.in
#!/sbin/openrc-run
# Copyright 1999-2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command=/usr/sbin/GLogiKd
name="GLogiK Daemon"
description="Daemon to handle special features on gaming keyboard"
description_start="Start the GLogiK daemon"
description_stop="Stop the GLogiK daemon"

# The special RC_SVCNAME variable contains the name of this service.
pidfile="/run/${RC_SVCNAME}.pid"

command_args="-d -p ${pidfile}"

depend() {
	use logger
	after bootmisc
	need localmount @INIT_NEED_DBUS@
}
