# data/udev/99-GLogiK.rules.in
# do not edit this file, it will be overwritten on update
#
# udev rules file for GLogiK devices


# load uinput module if necessary
SUBSYSTEM=="input", RUN+="/sbin/modprobe uinput"
KERNEL=="uinput", GROUP="@GLOGIKD_GROUP@", MODE="@DEVICES_MODE@"

ACTION!="add", GOTO="glogik_rules_end"

SUBSYSTEM=="usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c22d", GROUP="@GLOGIKD_GROUP@", MODE="@DEVICES_MODE@"
SUBSYSTEM=="usb", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c22e", GROUP="@GLOGIKD_GROUP@", MODE="@DEVICES_MODE@"

LABEL="glogik_rules_end"
